# systemd/vpn-watchdog.service
[Unit]
Description=ProtonPi VPN health watchdog

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/vpn-watchdog
